# deploy/namr.service
[Unit]
Description=namr — Strava activity auto-renamer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/namr
EnvironmentFile=/opt/namr/.env
ExecStart=/opt/namr/.venv/bin/namr run
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/namr/data

[Install]
WantedBy=multi-user.target
